# systemd/renew-ssl.timer
[Unit]
Description=Monthly Renew SSL Cert

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target

