# systemd/renew-ssl.service
[Unit]
Description=Renew SSL Certs

[Service]
ExecStart=/root/acme.sh/acme.sh --renew-all

[Install]
WantedBy=multi-user.target

