# systemd/renew-ssl.service
[Unit]
Description=Renew SSL Certs

[Service]
ExecStart=/root/renew-ssl.sh

[Install]
WantedBy=multi-user.target

